# supervisor/syslogd.conf
[program:vncserver]
command=/usr/sbin/syslogd --no-detach --no-forward
autostart=true
autorestart=true
startsecs=5

